# cleanup-dst.service
[Unit]
Description=DST Violation Query Cleanup (zombie tabs + stale files)

[Service]
Type=oneshot
Environment="PATH=/home/openclaw/.npm-global/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/aiext/bin/python3 /home/openclaw/.claude/skills/DST违章查询/cleanup_daemon.py \
    --project-root /home/openclaw
